# deploy/supervisord.conf
[supervisord]
logfile=/var/log/quant-platform/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false

[program:quant-streamlit]
command=/home/%(ENV_USER)s/projects/quant-platform/.venv/bin/streamlit run app.py --server.port=8501 --server.headless=true --server.address=0.0.0.0
directory=/home/%(ENV_USER)s/projects/quant-platform
user=%(ENV_USER)s
autostart=true
autorestart=true
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/quant-platform/streamlit.log
stdout_logfile_maxbytes=20MB
stderr_logfile=/var/log/quant-platform/streamlit_err.log
stderr_logfile_maxbytes=5MB
environment=PYTHONUNBUFFERED=1,LOG_LEVEL=INFO

[program:quant-alerts]
command=/home/%(ENV_USER)s/projects/quant-platform/.venv/bin/python -m scheduler.alerts
directory=/home/%(ENV_USER)s/projects/quant-platform
user=%(ENV_USER)s
autostart=true
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/var/log/quant-platform/alerts.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/quant-platform/alerts_err.log
environment=PYTHONUNBUFFERED=1

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
